# etc/systemd/system/btrfs-snap@.timer
[Unit]
Description=BTRFS Snapshot Service Timer (%I)
PartOf=%i.mount

[Timer]
OnUnitActiveSec=5min
OnActiveSec=3min
Unit=btrfs-snap@%i.service
